# dashboard/supervisor.conf
[program:codepuma]
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
user=alih
directory=/home/alih/code-dot-org/dashboard
command=/home/alih/code-dot-org/dashboard/edraak.sh
