# cookbooks/cdo-newrelic/templates/default/newrelic_init.sh.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          cdo-newrelic
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: New Relic alerting
# Description:       Manage New Relic alert policies on startup and shutdown.
### END INIT INFO

<%= @env ? @env.map{|k,v|"export #{k}=#{v}"}.join("\n") : '' %>
NEWRELIC_RB=<%=@newrelic_rb%>
STATUS_FILE=/run/cdo-newrelic

case "$1" in
  start|restart|force-reload)
    /opt/chef/embedded/bin/ruby -r${NEWRELIC_RB} -e "puts NewRelicClient.new.enable_alerts('$(hostname)')" \
      && touch ${STATUS_FILE}
    exit $?
  ;;
  stop)
    /opt/chef/embedded/bin/ruby -r${NEWRELIC_RB} -e "puts NewRelicClient.new.disable_alerts('$(hostname)')" \
      && rm -f ${STATUS_FILE}
    exit $?
  ;;
  status)
    [ -a "${STATUS_FILE}" ]
    exit $?
  ;;
  *)
    echo "Usage: /etc/init.d/cdo-newrelic {start|stop|restart|force-reload}"
    exit 1
esac
