# cookbooks/cdo-apps/templates/default/puma.service.erb
[Unit]
Description=<%= @app_name %> service
After=network.target

# Configuration options copied from file generated by systemd-sysv-generator,
# then updated to match the recommendations at:
# https://github.com/puma/puma/blob/v5.6.5/docs/systemd.md#service-configuration
[Service]
Type=simple
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/<%= @app_name %> start
ExecStop=/etc/init.d/<%= @app_name %> stop

[Install]
WantedBy=multi-user.target