# cookbooks/cdo-apps/templates/default/puma.service.erb
[Unit]
Description=<%= @app_name %> service
Requires=<%=@app_name%>.socket

# Configuration options copied from file generated by systemd-sysv-generator
[Service]
Type=simple
User=<%=@user%>
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
Environment=<%=@export_env.map{|k,v| %Q{"#{k}=#{v}"}}.join(' ')%>
WorkingDirectory=<%=@app_root%>
ExecStart=/usr/local/bin/bundle exec --keep-file-descriptors puma -C <%=@src_file%> -e <%=@env%>
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target
