# cookbooks/cdo-apps/templates/default/puma.service.erb
[Unit]
Description=<%= @app_name %> service
After=network.target

# Configuration options represent a combination of Puma's recommendations at
# https://github.com/puma/puma/blob/v5.6.5/docs/systemd.md#service-configuration
# and our own existing configuration at
# https://github.com/code-dot-org/code-dot-org/blob/e38737f7befb317871dbf18f46c8f65b3dc282ab/cookbooks/cdo-apps/templates/default/puma.sh.erb
[Service]
Type=simple

User=<%= @user %>
WorkingDirectory=<%= @app_root %>

# Vestigial piece of old infrastructure here; our sockets are still configured
# to get created at /run/unicorn from back when we used that
RuntimeDirectory=unicorn

<%= @export_env ? @export_env.map{|k,v|"Environment=#{k}=#{v}"}.join("\n") : '' %>
Environment=LANG=en_US.UTF-8

ExecStart=/usr/bin/bundle exec puma -C <%= @src_file %> -e <%= @env %>
Restart=always

[Install]
WantedBy=multi-user.target
